# com.groutledge.crate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    LaunchAgent: com.groutledge.crate
    Starts the crate singleton container at login and keeps it running.

    This file is the SOURCE template. It is installed into ~/Library/LaunchAgents
    by `scripts/deploy`, which also materializes the runtime artifacts it points
    at (~/.local/opt/crate, ~/.local/bin/docker-bin). The agent NEVER reads the
    git repo (~/dev/crate) — only the deployed copies under ~/.local.

    Install/reload is handled by `scripts/deploy`. Manual uninstall:
        launchctl unload ~/Library/LaunchAgents/com.groutledge.crate.plist
        rm ~/Library/LaunchAgents/com.groutledge.crate.plist

    The container uses restart: unless-stopped in compose.yaml, so Docker itself
    will restart it if it crashes. launchd only needs to run `docker compose up`
    once at login.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.groutledge.crate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Points at DEPLOYED runtime artifacts under ~/.local, never the repo.
             Resolve docker via the deployed helper (avoids dead OrbStack symlink
             at /usr/local/bin/docker; handles launchd's minimal PATH). Then wait
             for the daemon (Colima) and bring the container up. up -d is
             idempotent: no-op if already running. -->
        <string>
            DOCKER="$(/Users/groutledge/.local/bin/docker-bin)" || exit 1;
            until "$DOCKER" info &gt;/dev/null 2&gt;&amp;1; do sleep 2; done;
            cd /Users/groutledge/.local/opt/crate &amp;&amp; "$DOCKER" compose up -d
        </string>
    </array>

    <!-- Run at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep launchd from restarting this job after it exits successfully.
         docker compose up -d exits 0 once the container is running; the
         container itself stays alive via Docker's restart policy. -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/crate-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/crate-launchd.log</string>
</dict>
</plist>
