# com.groutledge.colima.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    LaunchAgent: com.groutledge.colima
    Starts the Colima VM at login — the foundation crate's and vault's own
    launchd jobs wait on (via `docker info`) before bringing their containers up.

    Boot chain, bottom to top:
        login -> colima start                          (this plist)
               -> crate:  wait for docker info, then compose up -d
               -> vault:  wait for docker info, then compose up -d

    launchd gives no ordering guarantee between separate LaunchAgents, so
    crate's and vault's plists don't depend on this one having already fired —
    they just poll `docker info` until it succeeds. That works whether this
    plist runs first, last, or Colima was already running from a prior session.

    `colima start` on an already-running instance exits fast (~1s) as a no-op
    — safe to invoke unconditionally on every login, same idempotent pattern
    used for `docker network create crate-net` in the other plists.

    ASSUMPTION: bare `colima start` resumes the default profile with whatever
    resources (CPU/memory/disk) were configured the last time it was started
    manually. If you use a named profile via colima's "profile" flag, add
    that flag and its value to the ProgramArguments string below.

    Install/reload is handled by crate's `scripts/deploy`. Manual uninstall:
        launchctl unload ~/Library/LaunchAgents/com.groutledge.colima.plist
        rm ~/Library/LaunchAgents/com.groutledge.colima.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.groutledge.colima</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            COLIMA="$(/Users/groutledge/.local/bin/colima-bin)" || exit 1;
            "$COLIMA" start
        </string>
    </array>

    <!-- Run at login. On a single-user Mac without a login delay (e.g.
         FileVault + auto-login), this is functionally "at boot". If you want
         it to run before any login (e.g. an always-on Mac mini), that needs
         a LaunchDaemon in /Library/LaunchDaemons instead — different
         installation path, runs as root, ask if you want that variant. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- colima start exits 0 once the VM is up; Colima itself stays running
         as its own background process — nothing here needs to be kept alive. -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/colima-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/colima-launchd.log</string>
</dict>
</plist>
