# docker/com.groutledge.crate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    LaunchAgent: com.groutledge.crate
    Starts the crate singleton container at login and keeps it running.

    This file is the SOURCE template. `@CRATE_HOME@` is a placeholder — launchd
    plists can't expand $HOME or ~ themselves, so `docker/deploy` substitutes it
    with the real path (from CRATE_HOME in docker/.env, default $HOME) when it
    renders this into ~/Library/LaunchAgents. The checked-in file is never
    edited directly; only the rendered copy has a real path. It also
    materializes the runtime artifacts this plist points at (~/.local/opt/crate,
    ~/.local/bin/docker-bin). The agent NEVER reads the git repo (~/dev/crate)
    — only the deployed copies under ~/.local.

    Install/reload is handled by `docker/deploy`. Manual uninstall:
        launchctl unload ~/Library/LaunchAgents/com.groutledge.crate.plist
        rm ~/Library/LaunchAgents/com.groutledge.crate.plist

    The container uses restart: unless-stopped in compose.yaml, so Docker itself
    will restart it if it crashes. launchd only needs to run `docker compose up`
    once at login.

    crate joins the `crate-net` Docker network, shared with the separate
    `vault` container (see ~/dev/vault). Both this plist and vault's own
    plist idempotently ensure the network exists before bringing their stack
    up, so boot order between the two launchd jobs never matters — whichever
    runs first creates the network, the other is a no-op.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.groutledge.crate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Points at DEPLOYED runtime artifacts under ~/.local, never the repo.
             Resolve docker via the deployed helper (avoids dead OrbStack symlink
             at /usr/local/bin/docker; handles launchd's minimal PATH). Then wait
             for the daemon (Colima) and bring the container up. up -d is
             idempotent: no-op if already running. -->
        <string>
            DOCKER="$(@CRATE_HOME@/.local/bin/docker-bin)" || exit 1;
            until "$DOCKER" info &gt;/dev/null 2&gt;&amp;1; do sleep 2; done;
            "$DOCKER" network create crate-net &gt;/dev/null 2&gt;&amp;1 || true;
            cd @CRATE_HOME@/.local/opt/crate &amp;&amp; "$DOCKER" compose up -d
        </string>
    </array>

    <!-- Run at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep launchd from restarting this job after it exits successfully.
         docker compose up -d exits 0 once the container is running; the
         container itself stays alive via Docker's restart policy. -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/crate-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/crate-launchd.log</string>
</dict>
</plist>
